# db/image-files/var/spool/cron/crontabs/root
# ### Test Cron Running
*   *   *   *   *   echo "Test Cron is Running" >> /var/log/cron.log 2>&1

